# home/.config/systemd/user/url-open@.service
[Unit]
Description=URL open handler (per-connection)

[Service]
Type=oneshot
ExecStart=/home/nathan.fritz/.scripts/url-enqueue
StandardInput=socket
